# systemd/oracle-report.service.template
[Unit]
Description=Oracle on-device Flask app
After=network-online.target llama-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{ORACLE_DIR}}
EnvironmentFile=-{{ORACLE_DIR}}/.env
ExecStart={{ORACLE_DIR}}/run.sh
Restart=on-failure
RestartSec=5
User={{ORACLE_USER}}

[Install]
WantedBy=multi-user.target
